# etc/profile.d/alt.sh
export PATH="$HOME/.local/alt/shims:$PATH"
