# etc/profile.d/alt.sh
# shellcheck shell=sh
export PATH="$HOME/.local/alt/shims:$PATH"
